# ext/templates/puppetdb.service.erb
[Unit]
Description=Puppet Enterprise <%= @name %> Daemon
After=syslog.target network.target

[Service]
Type=simple
EnvironmentFile=/etc/sysconfig/<%= @name %>
User=<%= @name %>
PIDFile=/var/run/<%= @name %>/<%= @name %>.pid
ExecStart=<%= @java_bin %> \
          $JAVA_ARGS \
          -jar ${INSTALL_DIR}/puppetdb.jar \
          services -c ${CONFIG} $@
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
